# .config/systemd/user/sBatteryctl.service
[Unit]
Description="sBatteryctl battery monitor"

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sBatteryctl monitor
